# deploy/systemd/kyuubiki-agent-registry.service
[Unit]
Description=Kyuubiki Agent Registry Sync
After=network-online.target kyuubiki-orchestrator.service kyuubiki-agent.service
Wants=network-online.target

[Service]
Type=simple
User={{KYUUBIKI_SERVICE_USER}}
Environment=KYUUBIKI_HOME={{KYUUBIKI_INSTALL_ROOT}}
ExecStart=/bin/sh -lc 'cd "$KYUUBIKI_HOME"; exec ./deploy/kyuubiki-agent-registry-sync.sh'
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
