# deploy/systemd/kyuubiki-deploy-server.service
[Unit]
Description=Kyuubiki Deploy Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{KYUUBIKI_SERVICE_USER}}
Environment=KYUUBIKI_HOME={{KYUUBIKI_INSTALL_ROOT}}
Environment=RUST_LOG=info
WorkingDirectory={{KYUUBIKI_INSTALL_ROOT}}
ExecStart={{KYUUBIKI_INSTALL_ROOT}}/workers/rust/target/debug/kyuubiki-deploy-server serve --host 127.0.0.1 --port 4070
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
